# ops/systemd/atlas-live-shadow.service
[Unit]
Description=Atlas live/shadow daily loop (forge->live; shadow on live data, $0)
After=network-online.target
[Service]
Type=oneshot
WorkingDirectory=/root/atlas
ExecStart=/usr/bin/python3 -m live.daily --mode shadow
StandardOutput=append:/tmp/atlas-live-shadow.log
StandardError=append:/tmp/atlas-live-shadow.log
TimeoutStartSec=1800
[Install]
WantedBy=multi-user.target
