# systemd/atlas-research-window.service
[Unit]
Description=Atlas Research Window — time-boxed parameter sweep
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/atlas
ExecStart=/root/atlas/scripts/research_cron.sh
EnvironmentFile=-/etc/atlas/atlas.conf
# 1h/worker × 5 strategies ÷ 3 workers ≈ 2h sweep + 25m LLM + buffer
TimeoutStartSec=9900
# Give SIGTERM 60s to clean up pool workers, then SIGKILL everything
TimeoutStopSec=60
# Kill ALL processes in the cgroup at once (prevents orphaned workers)
KillMode=control-group
SendSIGKILL=yes
# Nice priority so trading crons still preempt
Nice=10
# Low OOM priority (protect trading services)
OOMScoreAdjust=500
EnvironmentFile=/etc/openclaw/openclaw.env
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/root
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=TZ=Australia/Brisbane
