# systemd/atlas-reconcile-shadow.service
[Unit]
Description=Atlas reconcile shadow runner — B.2 cutover validation
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/root/atlas/data/HALT

[Service]
Type=oneshot
WorkingDirectory=/root/atlas
EnvironmentFile=-/etc/atlas/atlas.conf
ExecStart=/usr/bin/flock -n /tmp/reconcile_shadow.lock /usr/bin/timeout 300 /usr/bin/python3 scripts/reconcile_shadow.py --once --no-alert
StandardOutput=append:/root/atlas/logs/reconcile_shadow.log
StandardError=append:/root/atlas/logs/reconcile_shadow.log
TimeoutStartSec=360
