# systemd/atlas-director.timer
[Unit]
Description=Atlas Director — weekly queue review

[Timer]
# Sunday 08:00 AEST (Saturday 22:00 UTC)
OnCalendar=Sat 22:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
