# systemd/atlas-orchestrator.service
[Unit]
Description=Atlas Orchestrator (Phase C.3 — shadow mode)
Documentation=file:///root/atlas/docs/phase-c-orchestrator-timer.md
After=network.target

[Service]
Type=oneshot
User=root
EnvironmentFile=-/etc/atlas/atlas.conf
WorkingDirectory=/root/atlas
ExecStart=/usr/bin/python3 -m core.orchestrator --once --shadow
StandardOutput=append:/root/atlas/logs/orchestrator.log
StandardError=append:/root/atlas/logs/orchestrator.log
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
