# systemd/atlas-orchestrator.timer
[Unit]
Description=Atlas Orchestrator (15-min cycle)

[Timer]
OnCalendar=*:0/15
Persistent=false
Unit=atlas-orchestrator.service

[Install]
WantedBy=timers.target
