# systemd/atlas-director.service
[Unit]
Description=Atlas Director — automated queue management and portfolio review
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/atlas
ExecStart=/usr/bin/python3 -u scripts/director_cron.py
EnvironmentFile=-/etc/atlas/atlas.conf
TimeoutStartSec=4200
Nice=10
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/root
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=TZ=Australia/Brisbane
