# systemd/atlas-research-window@gold_etfs.timer
[Unit]
Description=Atlas Research Window — gold_etfs nightly sweep

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
Unit=atlas-research-window@gold_etfs.service

[Install]
WantedBy=timers.target
