# systemd/atlas-intraday-backfill.timer
# ─────────────────────────────────────────────────────────────────────────────
# atlas-intraday-backfill.timer
#
# Fires atlas-intraday-backfill.service daily at 23:30 UTC.
# That is ~30 min after US market close (16:00 ET = 20:00 UTC) + some
# settlement/data-propagation buffer.
#
# DISABLED BY DEFAULT — enable with:
#   systemctl enable atlas-intraday-backfill.timer
#   systemctl start  atlas-intraday-backfill.timer
# after operator approves the storage plan.
# ─────────────────────────────────────────────────────────────────────────────
[Unit]
Description=Atlas Intraday 5-min Backfill Timer (nightly 23:30 UTC)
Documentation=file:///root/atlas/docs/decisions/2026-05-18-intraday-5min-backfill.md

[Timer]
# Fire at 23:30 UTC every day
OnCalendar=*-*-* 23:30:00

# If the system was down when it should have fired, run it within 5 min of boot
Persistent=true
RandomizedDelaySec=120

# Unit to activate (default: same name with .service suffix)
Unit=atlas-intraday-backfill.service

[Install]
WantedBy=timers.target
