# systemd/atlas-discovery.timer
[Unit]
Description=Atlas Discovery — weekly paper discovery

[Timer]
# Wednesday 10:00 AEST (Wednesday 00:00 UTC) — midweek, non-competing with director
OnCalendar=Wed 00:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
