# systemd/atlas-research-window.service
[Unit]
Description=Atlas Research Window — time-boxed parameter sweep
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/atlas
ExecStart=/root/atlas/scripts/research_cron.sh
# 2h default + 5min buffer
TimeoutStartSec=7500
# Nice priority so trading crons still preempt
Nice=10
# Low OOM priority (protect trading services)
OOMScoreAdjust=500
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/root
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=TZ=Australia/Brisbane
