# systemd/atlas-telegram-bot.service
[Unit]
Description=Atlas Telegram Bot
After=network.target docker.service

[Service]
Type=simple
WorkingDirectory=/root/atlas
ExecStart=/usr/bin/python3 /root/atlas/services/telegram_bot.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/root

[Install]
WantedBy=multi-user.target
