# systemd/atlas-dashboard.service
[Unit]
Description=Atlas Dashboard Server (FastAPI)
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/atlas
ExecStart=/usr/bin/python3 -m uvicorn atlas.dashboard.app:app --host 127.0.0.1 --port 8899
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/root/atlas

StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-dashboard

[Install]
WantedBy=multi-user.target
