# systemd/atlas-heartbeat-watchdog.service
[Unit]
Description=Atlas heartbeat watchdog
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/atlas
ExecStart=/usr/bin/python3 /root/atlas/scripts/heartbeat_watchdog.py
EnvironmentFile=-/etc/atlas/atlas.conf
StandardOutput=journal
StandardError=journal
