# systemd/atlas-research-runner.service
[Unit]
Description=Atlas Research Runner — queue-based experiment execution daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/atlas
ExecStart=/usr/bin/python3 -u scripts/research_runner.py
Restart=on-failure
RestartSec=60
# Low priority so sweep and trading crons preempt
Nice=12
OOMScoreAdjust=500
# Graceful shutdown: SIGTERM → 30s → SIGKILL
TimeoutStopSec=30
KillMode=mixed
SendSIGKILL=yes
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/root
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=TZ=Australia/Brisbane

[Install]
WantedBy=multi-user.target
