# systemd/atlas-sandbox-9strats.service
[Unit]
Description=Atlas Sandbox — 9 strategy parameter sweeps (sp500)

[Service]
Type=simple
WorkingDirectory=/root/atlas
EnvironmentFile=-/etc/atlas/atlas.conf
EnvironmentFile=-/etc/openclaw/openclaw.env
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/root
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=TZ=Australia/Brisbane
ExecStart=/bin/bash /root/atlas/scripts/sandbox_9_strategies.sh
TimeoutStartSec=54000
TimeoutStopSec=120
KillMode=control-group
Nice=10
StandardOutput=append:/root/atlas/logs/atlas-sandbox-9strats.log
StandardError=append:/root/atlas/logs/atlas-sandbox-9strats.log

[Install]
WantedBy=multi-user.target
