# systemd/atlas-research-window@sector_etfs.timer
[Unit]
Description=Atlas Research Window — sector_etfs nightly sweep

[Timer]
OnCalendar=*-*-* 01:00:00
Persistent=true
Unit=atlas-research-window@sector_etfs.service

[Install]
WantedBy=timers.target
