# systemd/atlas-dashboard.service
[Unit]
Description=Atlas Dashboard Server (auth-protected)
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/atlas
ExecStart=/usr/bin/python3 /root/atlas/services/dashboard_server.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/root/atlas

StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-dashboard

[Install]
WantedBy=multi-user.target
