# systemd/atlas-intraday-backfill.service
# ─────────────────────────────────────────────────────────────────────────────
# atlas-intraday-backfill.service
#
# Tiingo 5-minute intraday OHLCV backfill for Atlas trading system.
# Incrementally pulls 5-min bars for all sp500 tickers via Tiingo IEX API.
#
# DISABLED BY DEFAULT — enable with:
#   systemctl enable atlas-intraday-backfill.timer
#   systemctl start  atlas-intraday-backfill.timer
# after operator approves the storage plan (see docs/decisions/2026-05-18-intraday-5min-backfill.md).
#
# ─────────────────────────────────────────────────────────────────────────────
[Unit]
Description=Atlas Intraday 5-min OHLCV Backfill (Tiingo IEX)
Documentation=file:///root/atlas/docs/decisions/2026-05-18-intraday-5min-backfill.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/atlas

# ─── Date calculation ────────────────────────────────────────────────────────
# --start = 90 days ago (rolling window for incremental updates)
# --end   = yesterday (avoid fetching today's incomplete bars)
# For initial full backfill, run manually with --start 2024-01-01.
ExecStart=/bin/bash -c '\
  START=$(date -d "90 days ago" +%%Y-%%m-%%d) && \
  END=$(date -d "yesterday"    +%%Y-%%m-%%d) && \
  python3 -m scripts.backfill_intraday_5min \
    --universe sp500 \
    --start "$START" \
    --end   "$END" \
'

# ─── Runtime limits ──────────────────────────────────────────────────────────
# 4-hour timeout: sufficient for 204-ticker × 3-month incremental update
# on a paid Tiingo plan (1.5s/call × 204 × 3 = ~15 min; 4h is headroom).
TimeoutStartSec=14400

# ─── Output & logging ────────────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-intraday-backfill

# ─── Environment ─────────────────────────────────────────────────────────────
# Set TIINGO_CALL_DELAY=72 for free-tier Tiingo (50 req/hr).
# Leave unset (defaults to 1.5s) for paid tier.
# EnvironmentFile=-/etc/atlas/atlas-intraday.conf

# ─── Security ────────────────────────────────────────────────────────────────
User=root
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
