# systemd/atlas-error-remediation.service
[Unit]
Description=Atlas autonomous error remediation — classifier cycle (Phase 3: AUTO_FIX active)
After=network-online.target
ConditionPathExists=!/root/atlas/data/HALT
ConditionPathExists=!/root/atlas/.live_halt
ConditionPathExists=!/root/atlas/data/AUTO_REMEDIATION_HALT

[Service]
Type=oneshot
WorkingDirectory=/root/atlas
Environment="TZ=Australia/Brisbane"
EnvironmentFile=-/etc/atlas/atlas.conf
EnvironmentFile=-/root/atlas/.env
ExecStart=/usr/bin/python3 /root/atlas/core/error_monitor.py --once --batch-size 50 --no-dry-run
TimeoutStartSec=300
StandardOutput=append:/root/atlas/logs/error-remediation.log
StandardError=append:/root/atlas/logs/error-remediation.log

[Install]
WantedBy=multi-user.target
