# systemd/atlas-research-window@sp500.timer
[Unit]
Description=Atlas Research Window — sp500 nightly sweep

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true
Unit=atlas-research-window@sp500.service

[Install]
WantedBy=timers.target
