# systemd/atlas-research-window.timer
[Unit]
Description=Atlas Research — scheduled sweep windows with idle gaps

[Timer]
# Weekdays: 5 windows × 2h = 10h/day, with ~1h cooldowns
# Window 1: 09:30  (after postclose settles)
# Window 2: 12:30  (after 1h cooldown)
# Window 3: 15:30  (after 1h cooldown, before healthz at 18)
# Window 4: 20:00  (after premarket, overnight quiet)
# Window 5: 23:00  (late night, 1h gap before intraday)
OnCalendar=Mon..Fri 09:30
OnCalendar=Mon..Fri 12:30
OnCalendar=Mon..Fri 15:30
OnCalendar=Mon..Fri 20:00
OnCalendar=Mon..Fri 23:00

# Weekends: 4 windows × 2h = 8h/day, no trading crons to avoid
OnCalendar=Sat,Sun 09:00
OnCalendar=Sat,Sun 12:00
OnCalendar=Sat,Sun 15:00
OnCalendar=Sat,Sun 20:00

# Total: ~50h weekday + ~16h weekend = ~66h/week research
Unit=atlas-research-window.service
Persistent=true

[Install]
WantedBy=timers.target
