# systemd/atlas-discovery.service
[Unit]
Description=Atlas Discovery — LLM paper-to-strategy pipeline
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/atlas
ExecStart=/usr/bin/python3 -u research/discovery/run.py
TimeoutStartSec=3600
Nice=15
EnvironmentFile=-/etc/openclaw/openclaw.env
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/root
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=TZ=Australia/Brisbane
