# systemd/atlas-research-window@defensive_etfs.timer
[Unit]
Description=Atlas Research Window — defensive_etfs nightly sweep

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
Unit=atlas-research-window@defensive_etfs.service

[Install]
WantedBy=timers.target
