# systemd/atlas-reconcile-shadow.timer
[Unit]
Description=Atlas reconcile shadow timer — every 30min during US RTH (Tue-Sat)

[Timer]
OnCalendar=Tue..Sat *-*-* 00,01,02,03,04,05,06,07:00,30 UTC
Persistent=false
Unit=atlas-reconcile-shadow.service

[Install]
WantedBy=timers.target
