# systemd/atlas-error-remediation.timer
[Unit]
Description=Atlas error-remediation classifier — 5min RTH / 15min off-hours

[Timer]
# 5-min cadence — system-wide.  The classifier itself decides whether to act
# (in market hours, returns ESCALATE_DEFERRED for unknowns; continues for
# IGNORE patterns).  Running the timer constantly keeps capture-fresh.
OnCalendar=*:0/5
RandomizedDelaySec=30
Persistent=false
Unit=atlas-error-remediation.service

[Install]
WantedBy=timers.target
