# systemd/atlas-fred-health.service
[Unit]
Description=Atlas FRED API health check
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/atlas
ExecStart=/usr/bin/flock -n /tmp/atlas_fred_health.lock /usr/bin/timeout 180 /usr/bin/python3 /root/atlas/scripts/check_fred_health.py
EnvironmentFile=-/etc/atlas/atlas.conf
StandardOutput=journal
StandardError=journal
