# deploy/cronus-watchdog.service
[Unit]
Description=Cronus Agent Watchdog
After=network.target docker.service
Documentation=https://github.com/cronus-trade/cronus

[Service]
Type=oneshot
WorkingDirectory=/root/cronus
Environment=CRONUS_CONFIG=/root/cronus/config/trading_paper.yaml
Environment=PYTHONPATH=/root/cronus
User=root
ExecStart=/usr/bin/python3 /root/cronus/scripts/watchdog.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
